# deploy/systemd/resumesystem-health.timer
[Unit]
Description=Run ResumeSystem health checks every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
RandomizedDelaySec=20s

[Install]
WantedBy=timers.target
